# scripts/nidas_eol_relay/eol_relay.crontab
PROJECT=MARSHALL2023

* * * * *    $HOME/isfs/scripts/nidas_eol_relay/check_nidas_udp_relay.sh start

